# build/rootfs/etc/gearbox/services/caddy/run
#!/bin/bash

/usr/local/bin/caddy -conf /usr/local/etc/caddy.conf

